# sources/crontab
* * * * * php /var/www/html/seat/artisan schedule:run